# core/moskv-1/com.moskv.v_omega.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moskv.v_omega</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/borjafernandezangulo/Documents/antigravity/lively-maxwell/scripts/v_omega_shield.py</string>
        <string>--kill</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/tmp/v_omega.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/v_omega.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
