# core/moskv-1/com.moskv.exergy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moskv.exergy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/sandbox-exec</string>
        <string>-f</string>
        <string>/Users/borjafernandezangulo/Documents/antigravity/lively-maxwell/core/moskv-1/kernel/moskv_sandbox.sb</string>
        <string>/bin/zsh</string>
        <string>/Users/borjafernandezangulo/.cortex/scripts/cron_exergy_monitor.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/moskv_exergy.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/moskv_exergy.log</string>
</dict>
</plist>
